# etc/init.d/network-manager-openrc
#!/sbin/openrc-run

description="Starts the NetworkManager daemon"
command="/usr/sbin/NetworkManager"
command_args="--no-daemon"
supervisor=supervise-daemon

NM_KEYFILE_DIR="/var/lib/NetworkManager"

depend() {
	need dbus udev
	provide net
}

start_pre() {
	checkpath -d -m 0755 /run/NetworkManager

	checkpath -d -m 0755 -o root:root "$NM_KEYFILE_DIR"

	find "$NM_KEYFILE_DIR" -maxdepth 1 -type f \( ! -user root -o ! -group root \) -print0 2>/dev/null \
		| xargs -0r chown root:root

	find "$NM_KEYFILE_DIR" -maxdepth 1 -type f -perm /077 -print0 2>/dev/null \
		| xargs -0r chmod 600
}
