# etc/init.d/network-manager-openrc
#!/sbin/openrc-run

description="Starts the NetworkManager daemon"
command="/usr/sbin/NetworkManager"
command_args="--no-daemon"
supervisor=supervise-daemon

depend() {
	need dbus udev
	provide net
}

start_pre() {
	checkpath -d -m 0755 /run/NetworkManager
}
